# assets/service/singbox-deploy-cert-renew.service.tmpl
[Unit]
Description=singbox-deploy certificate renewal
Wants=network-online.target
After=network-online.target

[Service]
Type=oneshot
User=root
ExecStart={{ .DeployBin }} cert renew --threshold-days {{ .ThresholdDays }}
