# template/service/singbox-deploy-monitor.service.tmpl
[Unit]
Description=singbox-deploy traffic monitor
After=network.target sing-box.service

[Service]
User=root
ExecStart={{ .DeployBin }} monitor serve --listen 127.0.0.1:{{ .MonitorPort }}
Restart=on-failure
RestartSec=10
MemoryMax=128M

[Install]
WantedBy=multi-user.target
